# table-tennis-scoreboard-dietpi/systemd/scoreboard-kiosk.service
[Unit]
Description=Table Tennis Scoreboard Kiosk Service
Documentation=https://github.com/nastynorm/table-tennis-scoreboard
After=graphical-session.target network-online.target
Wants=graphical-session.target network-online.target
Requires=network-online.target

[Service]
Type=simple
User=dietpi
Group=dietpi
Environment=DISPLAY=:0
Environment=HOME=/home/dietpi
Environment=XDG_RUNTIME_DIR=/run/user/1000
WorkingDirectory=/home/dietpi/scoreboard
ExecStartPre=/bin/bash -c 'until ping -c1 google.com; do sleep 1; done'
ExecStart=/home/dietpi/scripts/start-scoreboard.sh
ExecStop=/bin/bash -c 'pkill -f "serve dist"; pkill -f "chromium-browser.*kiosk"'
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStartSec=120
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/dietpi /tmp /var/log
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits for Pi Zero 2W
MemoryMax=300M
CPUQuota=80%

# Environment variables for Chromium
Environment=CHROMIUM_FLAGS="--no-sandbox --disable-gpu"

[Install]
WantedBy=graphical-session.target
Also=network-online.target